# build-pine64-image/rc.local
#!/bin/sh -e
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.

FLAG="/opt/afterfirstboot.lock"
LOG="/var/log/first-boot.log"
FAILED=0

# Only executed after first boot
if [ ! -f $FLAG ] && [ ! -f "/root/afterfirstboot.lock" ]; then #TODO: Remove 'root' after 2017-05
  (/bin/bash /boot/first-boot.sh)
fi

exit 0

# vim: filetype=sh
