# build-pine64-image/rc.local
#!/bin/sh -e
#
# This script is executed at the end of each multiuser runlevel.
# Make sure that the script will "exit 0" on success or any other
# value on error.

FLAG="/root/afterfirstboot.lock"
LOG="/var/log/raspbian-ua-netinst.log"
FAILED=0

# Only executed after first boot
if [ ! -f $FLAG ]; then
  echo "" >> $LOG
  echo "[openHABian] Booting for the first time!" >> $LOG
exit 0

# vim: filetype=sh
