# includes/generic/lc.timer
[Unit]
Description=Run l check weekly

[Timer]
Unit=lc.service
OnCalendar=weekly
OnBootSec=15min
Persistent=true

[Install]
WantedBy=timers.target
