# includes/firemotd.service
[Unit]
Description=Periodic FireMotD stats regeneration
After=network.target network-online.target

[Service]
Type=oneshot
User=root
Group=root
RemainAfterExit=no
ExecStart=/usr/local/bin/FireMotD -S -D all &> /dev/null

[Install]
WantedBy=multi-user.target
