# includes/grott.service
[Unit]
Description=Grott Proxy Service
After=network.target

[Service]
ExecStart=/usr/bin/python3 -u %DIR/grott.py -v
WorkingDirectory=%DIR
Restart=on-failure
User=%USR

[Install]
WantedBy=multi-user.target
