# includes/homegear.service
[Unit]
Description=Homegear

[Service]
Type=simple
User=homegear
Group=homegear
ExecStart=/usr/bin/homegear

[Install]
WantedBy=multi-user.target
