# includes/sdrawcopy.timer
[Unit]
Description=Run SD raw dump semiannually
Requires=sdrawcopy.service

[Timer]
Unit=sdrawcopy.service
OnCalendar=*-01,07-01 01:15:00
Persistent=true

[Install]
WantedBy=timers.target
